# service/nginx/finish
#!/bin/sh

exec kill `cat /web/run/nginx.pid`
