# rootfs/etc/service/nginx/run
#!/bin/sh
/usr/sbin/nginx -c /web/config/nginx.conf  -g "daemon off;pid /web/run/nginx.pid;" -p /web

